# signaged.service
[Unit]
Description=Digital Signage Daemon

[Service]
ExecStart=/usr/bin/signaged

[Install]
WantedBy=multi-user.target
